# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:qdrant]
command=/app/bin/qdrant
directory=/app/data/qdrant
environment=QDRANT__STORAGE__STORAGE_PATH="/app/data/qdrant",QDRANT__SERVICE__HTTP_PORT="6333",QDRANT__SERVICE__GRPC_PORT="6334",QDRANT__LOG_LEVEL="WARN"
autostart=true
autorestart=true
priority=10
startsecs=5
stopwaitsecs=15
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:embedding]
command=/app/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8001 --workers 1
directory=/app/python/embedding_server
environment=HF_HOME="/app/models",HF_HUB_OFFLINE="1",TRANSFORMERS_OFFLINE="1",EMBEDDING_MODEL="BAAI/bge-m3",EMBEDDING_CACHE_DIR="/app/data/embedding-cache"
autostart=true
autorestart=true
priority=20
startsecs=10
stopwaitsecs=15
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:reranker]
command=/app/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8002 --workers 1
directory=/app/python/reranker_server
environment=HF_HOME="/app/models",HF_HUB_OFFLINE="1",TRANSFORMERS_OFFLINE="1",RERANKER_MODEL="BAAI/bge-reranker-v2-m3"
autostart=true
autorestart=true
priority=20
startsecs=10
stopwaitsecs=15
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:api]
command=/app/docker/wait-for-services.sh node /app/dist/api/server.js
environment=PORT="3000",HOST="0.0.0.0",RAG_CONFIG="/app/rag.config.json",RAG_DATA_DIR="/app/data",NODE_ENV="production"
autostart=true
autorestart=true
priority=30
startsecs=5
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
